# deploy/com.hunter.nyc-apt-radar.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hunter.nyc-apt-radar</string>
  <key>ProgramArguments</key>
  <array>
  <string>/usr/bin/env</string>
  <string>PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  <string>npm</string>
  <string>run</string>
  <string>agent:run</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/hunter/Developer/nyc-apt-radar</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/Users/hunter/Developer/nyc-apt-radar/data/logs/agent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hunter/Developer/nyc-apt-radar/data/logs/agent.err.log</string>
</dict>
</plist>
