# deploy/nyc-apt-radar.timer.example
[Unit]
Description=Run NYC Apt Radar every 60 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=60min
AccuracySec=1min
Persistent=true
Unit=nyc-apt-radar.service

[Install]
WantedBy=timers.target
